# no.havoyet.flask.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>no.havoyet.flask</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/eriko/Desktop/Havøyet AS/Claud Cowork/Bestillingsside/app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/eriko/Desktop/Havøyet AS/Claud Cowork/Bestillingsside</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/eriko/Desktop/Havøyet AS/Claud Cowork/Bestillingsside/flask.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eriko/Desktop/Havøyet AS/Claud Cowork/Bestillingsside/flask.log</string>
</dict>
</plist>
